# Proyecto laboratorio/linaro-jessie-developer-20161117-32/binary/var/lib/dpkg/info/96boards-tools-common.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/resize-helper.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/70-96boards-common.rules
/etc
/etc/profile.d
/etc/profile.d/96boards-sensors.sh
/etc/kernel
/etc/kernel/postinst.d
/etc/kernel/postinst.d/zz-copy-devicetree
/etc/kernel/postrm.d
/etc/kernel/postrm.d/zz-copy-devicetree
/usr
/usr/sbin
/usr/sbin/resize-helper
/usr/sbin/grub-install-hikey
/usr/share
/usr/share/96boards-tools
/usr/share/96boards-tools/grub
/usr/share/96boards-tools/grub/grub.cfg.in
/usr/share/96boards-tools/grub/grub-fastboot.cfg
/usr/share/96boards-tools/grub/grub-sd.configfile
/usr/share/96boards-tools/grub/grub.configfile
/usr/share/96boards-tools/grub/grub-hikey.cfg.in
/usr/share/96boards-tools/flash-hikey
/usr/share/doc
/usr/share/doc/96boards-tools-common
/usr/share/doc/96boards-tools-common/changelog.gz
/usr/share/doc/96boards-tools-common/copyright
