# Proyecto laboratorio/linaro-jessie-developer-20161117-32/binary/var/lib/dpkg/info/netperf.conffiles
/etc/default/netperf
/etc/init.d/netperf
