# Proyecto laboratorio/linaro-jessie-developer-20161117-32/binary/var/lib/dpkg/info/netperf.list
/.
/usr
/usr/bin
/usr/bin/netserver
/usr/bin/netperf
/usr/share
/usr/share/doc
/usr/share/doc/netperf
/usr/share/doc/netperf/netperf.html
/usr/share/doc/netperf/changelog.Debian.gz
/usr/share/doc/netperf/Release_Notes.gz
/usr/share/doc/netperf/examples
/usr/share/doc/netperf/examples/udp_stream_script
/usr/share/doc/netperf/examples/snapshot_script.gz
/usr/share/doc/netperf/examples/packet_byte_script.gz
/usr/share/doc/netperf/examples/tcp_stream_script
/usr/share/doc/netperf/examples/sctp_stream_script
/usr/share/doc/netperf/examples/tcp_rr_script
/usr/share/doc/netperf/examples/tcp_range_script
/usr/share/doc/netperf/examples/udp_rr_script
/usr/share/doc/netperf/AUTHORS.gz
/usr/share/doc/netperf/README.Debian
/usr/share/doc/netperf/copyright
/usr/share/doc/netperf/netperf.pdf.gz
/usr/share/doc/netperf/README
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/netserver.1.gz
/usr/share/man/man1/netperf.1.gz
/etc
/etc/default
/etc/default/netperf
/etc/init.d
/etc/init.d/netperf
